# djangox/deploy/uwsgi.conf
# Emperor uWSGI script

description "uWSGI Emperor"
start on runlevel [2345]
stop on runlevel [06]

exec sudo -uubuntu /home/ubuntu/venv/bin/uwsgi --die-on-term --emperor /etc/uwsgi